# reader/crontab
SHELL=/bin/bash

# min hour day month dayofweek cmd
*/20 *     *   *     *         ~/releases/prod/run.sh >> ~/logs/stdout.log 2>&1
